# batesste-icloud-exporter.socket
[Socket]
ListenFIFO=%t/batesste-icloud-exporter.stdin
Service=batesste-icloud-exporter.service
